# crontab
* * * * * node /app/scraper.js >> /var/log/scraper.log 2>&1
